# system_files/etc/profile.d/update.sh
alias rollback="run0 bootc rollback"
alias update="/usr/share/spacefin/scripts/update.sh"
