# system_files/etc/profile.d/update.sh
alias rollback="run0 pacman-bootc rollback"
alias pacman="/usr/libexec/pacman-wrapper"
alias update="/usr/share/spacefin/scripts/update.sh"
