# system_files/shared/etc/profile.d/common.sh
alias open="xdg-open &>/dev/null"
alias nano=micro
alias sudo=run0
alias zed=zeditor
alias spacefin-cli="just -f /usr/share/spacefin/Justfile"

alias neofetch=fastfetch
alias f-paleofetch="fastfetch -c /usr/share/fastfetch/presets/paleofetch.jsonc"
alias f-neofetch="fastfetch -c /usr/share/fastfetch/presets/neofetch.jsonc"
alias f-ci="fastfetch -c /usr/share/fastfetch/presets/ci.jsonc"
alias f-all="fastfetch -c /usr/share/fastfetch/presets/all.jsonc"
