# system_files/shared/etc/profile.d/common.sh
alias open="xdg-open &>/dev/null"
alias neofetch=fastfetch
alias spacefin-cli="just -f /usr/share/spacefin/.Justfile"
