# system_files/shared/etc/profile.d/update.sh
alias rollback="sudo bootc rollback"
alias update="just -f /usr/share/spacefin/Justfile update"
