# system_files/shared/etc/profile.d/update.sh
alias rollback="sudo bootc rollback"
alias update="/usr/share/spacefin/scripts/update.sh"
